# ansible/ansible-kafka-zookeeper-role/upstart/kafka.conf
# description "Kafka upstart config"
# author      "Jon Hunter"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn
respawn limit 5 60

exec {% for key, value in kafka_brokers.iteritems() %} /usr/local/kafka/bin/kafka-server-start.sh /usr/local/kafka/config/server-{{ key }}.properties & {% endfor %}
