# ansible/ansible-kafka-zookeeper-role/upstart/zookeeper.conf
# description "Zookeeper upstart config"
# author      "Jon Hunter"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn
respawn limit 5 60

exec /usr/local/zookeeper/bin/zkServer.sh start
